# deploy/upstart.tpl.conf
description "%(server_name)s"
start on runlevel [2345]
stop on runlevel [016]
respawn
console log
chdir %(install_dir)s/src
exec python server.py %(port)s
